# deployment/supervisor/horizon.conf
# Supervisor configuration for Laravel Horizon (Production)
# Copy this file to: /etc/supervisor/conf.d/horizon.conf
#
# After copying, run:
# sudo supervisorctl reread
# sudo supervisorctl update
# sudo supervisorctl start horizon

[program:horizon]
process_name=%(program_name)s
command=php /var/www/pro_hub_started/artisan horizon
autostart=true
autorestart=true
user=deployer
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/pro_hub_started/storage/logs/horizon.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600
stopasgroup=true
killasgroup=true
